# templates/aem_type.erb
#!/bin/bash
#
# chkconfig: 35 99 01
# description: Starts and stops aem  services.
#
#  This file is managed by Puppet.  DON"T EDIT FILE
#

# Source function library.
. /etc/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 1

SCRIPTNAME=$0
INSTANCENAME="<%= @instance_name %>"
PROG="AEM ${INSTANCENAME} instance."
STARTUPDIR="<%= scope['::adobe_em6::params::dir_aem_install'] %>/${INSTANCENAME}/crx-quickstart"
AEMUSER="<%= scope['::adobe_em6::params::aem_user'] %>"
LOCK="/var/lock/subsys/${SCRIPTNAME}"

start() {
  [ -f "${LOCK}" ] && exit 0
  touch "${LOCK}"
  log "$(date) - Starting: ${PROG}"
  su - ${AEMUSER} -c "${STARTUPDIR}/bin/start"
}

stop() {
  log "$(date) - Shutting down: ${PROG}"
  su - ${AEMUSER} -c "${STARTUPDIR}/bin/stop"
  wait
  #Let things settle a bit
  sleep 5
  rm -rf "${LOCK}"
}

log() {
  echo $1
  logger $1
}

# See how we were called.
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart|reload)
        stop
        sleep 30
        start
        ;;
  status)
        su - ${AEMUSER} -c "${STARTUPDIR}/bin/status"
        ;;
  *)
        echo $"Usage: ${SCRIPTNAME} {start|stop|restart|status}"
        exit 1
esac
